# build/resilio-sync.conf
[program:resilio-sync]
user = nobody
command = /bin/bash -c "mkdir -p /config/resilio/licence /config/resilio/logs && /usr/bin/rslsync --storage /config/resilio/ --webui.listen 0.0.0.0:8888 --nodaemon --log /config/resilio/logs/resilio.log"
umask = 000

